# rc.local
#!/bin/bash
cd /mnt/internal_storage/
./supervisor.sh &